# etc/systemd/system/nix-daemon.service.d/50-arch-install-mount-before.conf
[Unit]
Requires=-.mount nix.mount tmp.mount usr.mount
After=-.mount nix.mount tmp.mount usr.mount
