# etc/udev/rules.d/50-arch-install.rules
ACTION=="add", KERNEL=="sd*[!0-9]", SUBSYSTEM=="block", RUN+="/usr/bin/hdparm --security-freeze /dev/%k"
