# etc/udev/rules.d/50-arch-install.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="/usr/bin/hdparm --security-freeze /dev/%k"
